# epaper.service
# E-paper image server – run on Pi at boot.
# Install with: sudo ./enable-boot.sh   (from repo root; replaces YOUR_USER with current user)
# If it fails: journalctl -u epaper -n 40 --no-pager
# Optional override: /etc/systemd/system/epaper.service.d/override.conf (EPD_LISTEN_PORT=8080 etc.)

[Unit]
Description=E-paper image server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/epaper
SupplementaryGroups=spi gpio
Environment=PYTHONPATH=/home/YOUR_USER/epaper/lib/e-Paper/RaspberryPi_JetsonNano/python/lib
Environment=EPD_LISTEN_HOST=0.0.0.0
Environment=EPD_LISTEN_PORT=8080
ExecStart=/usr/bin/python3 /home/YOUR_USER/epaper/app/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
