# epaper.service
[Unit]
Description=e-Paper photo display server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# install-service.sh replaces SERVICE_USER with your login user
# Change to root only if SPI/GPIO access fails with your user
User=SERVICE_USER
WorkingDirectory=REPO_DIR
ExecStart=/usr/bin/python3 REPO_DIR/display_photo.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
