# epaper.service
# StellarEars e-ink display client – run on Pi host at boot.
# Copy: sudo cp epaper.service /etc/systemd/system/
#       sudo systemctl daemon-reload && sudo systemctl enable epaper && sudo systemctl start epaper
#
# Optional: override env (e.g. different host/port for push listener):
#   sudo mkdir -p /etc/systemd/system/epaper.service.d
#   echo -e '[Service]\nEnvironment="EPD_LISTEN_HOST=0.0.0.0"\nEnvironment="EPD_LISTEN_PORT=9090"' | sudo tee /etc/systemd/system/epaper.service.d/override.conf
#   sudo systemctl daemon-reload && sudo systemctl restart epaper

[Unit]
Description=StellarEars e-ink display client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/epaper
Environment=PYTHONPATH=/home/pi/epaper/lib/e-Paper/RaspberryPi_JetsonNano/python/lib
Environment=EPD_LISTEN_HOST=0.0.0.0
Environment=EPD_LISTEN_PORT=9090
ExecStart=/usr/bin/python3 /home/pi/epaper/app/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
